# deployments/ansible/reconforge.service.j2
[Unit]
Description=ReconForge Reconnaissance Service
Documentation=https://github.com/reconforge/reconforge
After=network.target

[Service]
Type=simple
User={{ user }}
Group={{ group }}
WorkingDirectory={{ install_dir }}
ExecStart=/usr/local/bin/reconforge -t {{ target | default('example.com') }} --deep
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10
LimitNOFILE=65536

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths={{ output_dir }} {{ install_dir }}/logs

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=reconforge

[Install]
WantedBy=multi-user.target
